# nestlabel.service
[Unit]
Description=NestLabel Annotation System
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/annotation
Environment=PATH=/home/ubuntu/annotation/venv/bin
EnvironmentFile=/home/ubuntu/annotation/.env
ExecStart=/home/ubuntu/annotation/venv/bin/uvicorn web.backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
